# Backend/beauty-crm-sms-worker.ini
[program:beauty-crm-sms-worker]
command=php /www/wwwroot/api.ziboxcrm.ir/BCRM/Backend/artisan queue:work database --queue=sms,high,default --sleep=1 --tries=2 --timeout=60 --memory=128
directory=/www/wwwroot/api.ziboxcrm.ir/BCRM/Backend/
autostart=true
autorestart=true
startsecs=3
startretries=3
exitcodes=0,2
loglevel=info
stdout_logfile=/www/server/panel/plugin/supervisor/log/beauty-crm-sms-worker.out.log
stderr_logfile=/www/server/panel/plugin/supervisor/log/beauty-crm-sms-worker.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
redirect_stderr=true
user=www
priority=999
numprocs=2
stopsignal=QUIT
process_name=%(program_name)s_%(process_num)02d